# ops/supervisor/moss-all-in-one-supervisord.conf
[supervisord]
nodaemon=true
logfile=/dev/null
pidfile=/tmp/supervisord.pid
childlogdir=/tmp
logfile_maxbytes=0

[unix_http_server]
file=/tmp/supervisor.sock
chmod=0700

[supervisorctl]
serverurl=unix:///tmp/supervisor.sock

[rpcinterface:supervisor]
supervisor.rpcinterface_factory=supervisor.rpcinterface:make_main_rpcinterface

[program:moss-gateway]
command=/opt/hermes/.venv/bin/hermes gateway run
directory=/agents/moss/private
user=hermes
autostart=true
autorestart=true
priority=10
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=HOME="/opt/data",HERMES_HOME="/opt/data",HERMES_KANBAN_HOME="/mnt/hermes-shared/kanban",PATH="/opt/hermes/.venv/bin:/opt/data/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[program:moss-dashboard]
command=/opt/hermes/.venv/bin/hermes dashboard --host 0.0.0.0 --port 9119 --no-open --insecure
directory=/agents/moss/private
user=hermes
autostart=true
autorestart=true
priority=20
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=HOME="/opt/data",HERMES_HOME="/opt/data",PATH="/opt/hermes/.venv/bin:/opt/data/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"

[program:moss-webui]
command=/opt/hermes/.venv/bin/python /opt/hermes-webui/server.py
directory=/opt/hermes-webui
user=hermes
autostart=true
autorestart=true
priority=30
startsecs=5
stdout_logfile=/dev/fd/1
stdout_logfile_maxbytes=0
stderr_logfile=/dev/fd/2
stderr_logfile_maxbytes=0
environment=HOME="/opt/data",HERMES_HOME="/opt/data",HERMES_WEBUI_AGENT_DIR="/opt/hermes",HERMES_WEBUI_HOST="0.0.0.0",HERMES_WEBUI_PORT="8787",HERMES_WEBUI_STATE_DIR="/opt/data/webui",HERMES_WEBUI_DEFAULT_WORKSPACE="/agents/moss/private",PATH="/opt/hermes/.venv/bin:/opt/data/.local/bin:/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin"
